# btsync/supervisord.conf
[program:btsync]
command=btsync --config /etc/btsync.cfg --nodaemon
autostart=false
redirect_stderr=true

[program:confd]
command=confd -node %(ENV_ETCD_HOST)s
redirect_stderr=true

